# backend/deploy/supervisord-workers.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:market_stream]
command=python -m app.workers.market_stream
directory=/app
autostart=true
autorestart=true
startsecs=1
startretries=100
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:brain_track_worker]
command=python -m app.workers.brain_track_worker
directory=/app
autostart=true
autorestart=true
startsecs=1
startretries=100
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
